# meow-bot.service
[Unit]
Description=Meow Bot Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/myMeowBot
Environment=PATH=/root/myMeowBot/venv/bin:/usr/bin:/usr/local/bin
EnvironmentFile=/root/myMeowBot/.env
ExecStart=/bin/bash -c "source /root/myMeowBot/venv/bin/activate && python3 /root/myMeowBot/bot.py"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Настройки безопасности
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/myMeowBot/logs

[Install]
WantedBy=multi-user.target 