# meow-bot.service
[Unit]
Description=Meow Bot Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/myMeowBot
Environment=PATH=/root/myMeowBot/venv/bin:/usr/bin:/usr/local/bin
EnvironmentFile=/root/myMeowBot/.env
ExecStart=/root/myMeowBot/venv/bin/python3 /root/myMeowBot/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 