# meow-bot.service
[Unit]
Description=Meow Bot Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/myMeowBot
Environment=PATH=/root/myMeowBot/venv/bin
ExecStart=/root/myMeowBot/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Настройки безопасности
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/myMeowBot/logs

[Install]
WantedBy=multi-user.target 